# ytmusic-backend.service
[Unit]
Description=YTMusic Backend
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/ytmusic-backend
EnvironmentFile=/opt/ytmusic-backend/.env
# One worker so in-memory stream URL cache is shared; threads for concurrent streams.
ExecStart=/opt/ytmusic-backend/venv/bin/gunicorn -w 1 -k gthread --threads 8 -b 127.0.0.1:5000 app:app
Restart=always
RestartSec=3
TimeoutStopSec=20

[Install]
WantedBy=multi-user.target
